# debian/softether.upstart
description "SoftEther VPN"

start on started networking
stop on stoppping networking

respawn
expect stop

script
	test -f /etc/default/softether && . /etc/default/softether
	case $SOFTETHER_MODE in
	(vpnbridge|vpnclient|vpnserver)
		exec /usr/sbin/softetherd "$SOFTETHER_MODE" ;;
	(*)
		exit 0 ;;
	esac
end script
